# ansible/roles/snapper/templates/snapper-cleanup.timer.j2
[Unit]
Description=Daily Cleanup of Snapper Snapshots (2am)
Documentation=man:snapper(8) man:snapper-configs(5)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=false

[Install]
WantedBy=timers.target
