# config/systemd/user/agenda-sync.service
[Unit]
Description=Sync Google Calendar agenda

[Service]
Type=oneshot
ExecStart=/home/erik/.bin/agenda-sync
