# config/systemd/user/photos-sync.service
[Unit]
Description=Sync photos from NAS

[Service]
Type=oneshot
ExecStart=%h/.bin/photos-sync
