# config/systemd/user/scans-sync.service
[Unit]
Description=Sync scans from Google Drive

[Service]
Type=oneshot
ExecStart=%h/.bin/scans-sync
