# config/systemd/user/scans-sync.timer
[Unit]
Description=Sync scans from Google Drive every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
