# config/systemd/user/system-health.service
[Unit]
Description=Collect system health snapshot

[Service]
Type=oneshot
ExecStart=/home/erik/.bin/system-health --snapshot
