# system/etc/systemd/system/sddm.service.d/keyd.conf
[Unit]
After=keyd.service
Wants=keyd.service
